# ScriptService.service
[Unit]
Description=ChromeOs Integration
 
[Service]
ExecStart=/bin/cros-upgrade
 
[Install]
WantedBy=default.target
